# config/validation-system.service
[Unit]
Description=Perfect Validation System
Documentation=https://docs.validation-system.com
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=simple
User=validation-system
Group=validation-system
WorkingDirectory=/opt/validation-system
ExecStart=/opt/validation-system/bin/validation-daemon --config /etc/validation-system/validation-system.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/validation-system /var/lib/validation-system
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=200%

# Environment
Environment=VALIDATION_SYSTEM_CONFIG=/etc/validation-system/validation-system.conf
Environment=VALIDATION_SYSTEM_LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
